# systemd/.config/systemd/user/set-wallpaper.service
[Unit]
Description=Set Random Wallpaper

[Service]
Type=oneshot
ExecStart=/bin/bash ~/.config/scripts/randomWallpaper.sh
