# systemd/.config/systemd/user/set-gnome-wallpaper.service
[Unit]
Description=Set random GNOME wallpaper

[Service]
Type=oneshot
ExecStart=/home/rodney/.config/scripts/set-gnome-wallpaper.sh

